# supervisor.conf
[supervisord]
nodaemon=true

[program:nginx]
command=/usr/sbin/nginx -g 'daemon off;'

[program:php5-fpm]
command=/usr/sbin/php-fpm -F

[program:crond]
command=/usr/sbin/crond

